# registry/prometheus/quadlets/prometheus.container
[Unit]
Description=Prometheus - time-series metrics
StartLimitIntervalSec=60
StartLimitBurst=3

[Container]
ContainerName=prometheus
Image=docker.io/prom/prometheus:v2.55.1
Network=prometheus.network
PublishPort=${PORT_HTTP}:9090
Volume=%h/.local/share/services/prometheus/configs/prometheus.yml:/etc/prometheus/prometheus.yml:Z,ro
Volume=%h/.local/share/services/prometheus/targets:/etc/prometheus/targets:Z,ro
Volume=%h/.local/share/services/prometheus/data:/prometheus:Z,U
Exec=--config.file=/etc/prometheus/prometheus.yml --storage.tsdb.path=/prometheus --web.console.libraries=/usr/share/prometheus/console_libraries --web.console.templates=/usr/share/prometheus/consoles --web.enable-lifecycle
Pull=missing
EnvironmentFile=%h/.local/share/services/prometheus/.env
HealthCmd=wget -q --spider http://localhost:9090/-/ready || exit 1

[Service]
EnvironmentFile=%h/.local/share/services/prometheus/.env
Restart=always
RestartSec=5
TimeoutStartSec=300

[Install]
WantedBy=default.target
