# registry/uptime-kuma/quadlets/uptime-kuma.container
[Unit]
Description=Uptime Kuma
StartLimitIntervalSec=60
StartLimitBurst=3

[Container]
ContainerName=uptime-kuma
Image=docker.io/louislam/uptime-kuma:1
PublishPort=${RYRA_PORT_HTTP}:3001
Volume=uptime-kuma-data.volume:/app/data:U
Pull=never
EnvironmentFile=%h/.local/share/ryra/uptime-kuma/.env

[Service]
EnvironmentFile=%h/.local/share/ryra/uptime-kuma/.env
Restart=always
RestartSec=5
TimeoutStartSec=300

[Install]
WantedBy=default.target
