# registry/zammad/quadlets/zammad.container
[Unit]
Description=Zammad nginx (web entry point)
Requires=zammad-railsserver.service zammad-websocket.service
After=zammad-railsserver.service zammad-websocket.service
StartLimitIntervalSec=60
StartLimitBurst=3

[Container]
Image=ghcr.io/zammad/zammad:7.0.1-0019
ContainerName=zammad
Network=zammad.network
PublishPort=${SERVICE_PORT_HTTP}:8080
Volume=%h/.local/share/services/zammad/storage:/opt/zammad/storage:Z,U
Exec=zammad-nginx
HealthCmd=curl -sf http://localhost:8080/ || exit 1
HealthStartPeriod=60s
HealthInterval=30s
HealthRetries=5
HealthTimeout=5s
Pull=missing
EnvironmentFile=%h/.local/share/services/zammad/.env

[Service]
EnvironmentFile=%h/.local/share/services/zammad/.env
Restart=always
RestartSec=5
TimeoutStartSec=600
ExecStartPost=/bin/bash %h/.local/share/services/zammad/configs/scripts/run-autowizard.sh
ExecStartPost=/bin/bash %h/.local/share/services/zammad/configs/scripts/configure-oidc.sh
ExecStartPost=/bin/bash %h/.local/share/services/zammad/configs/scripts/configure-smtp.sh

[Install]
WantedBy=default.target
