# registry/zammad/quadlets/zammad-elasticsearch.container
[Unit]
Description=Zammad Elasticsearch
StartLimitIntervalSec=60
StartLimitBurst=3

[Container]
Image=docker.io/library/elasticsearch:9.3.3
ContainerName=zammad-elasticsearch
Network=zammad.network:alias=elasticsearch
Volume=zammad-es-data.volume:/usr/share/elasticsearch/data:U
Environment=discovery.type=single-node
Environment=xpack.security.enabled=false
Environment=ES_JAVA_OPTS=-Xms512m -Xmx512m
HealthCmd=curl -sf http://localhost:9200/_cluster/health || exit 1
HealthStartPeriod=60s
HealthInterval=30s
HealthRetries=10
HealthTimeout=5s
Pull=missing

[Service]
Restart=always
RestartSec=5
TimeoutStartSec=600

[Install]
WantedBy=default.target
