# registry/caddy/quadlets/caddy.container
[Unit]
Description=Caddy reverse proxy
StartLimitIntervalSec=60
StartLimitBurst=3

[Container]
ContainerName=caddy
Image=docker.io/library/caddy:2-alpine
AutoUpdate=registry
Network=caddy.network
PublishPort=${PORT_HTTP}:${PORT_HTTP}
PublishPort=${PORT_HTTPS}:${PORT_HTTPS}
Volume=caddy-data.volume:/data:U
Volume=caddy-config-state.volume:/config:U
Volume=%h/services/caddy/config:/etc/caddy:Z
Pull=missing
EnvironmentFile=%h/services/caddy/.env
HealthCmd=caddy version || exit 1

[Service]
EnvironmentFile=%h/services/caddy/.env
Restart=always
RestartSec=5
TimeoutStartSec=300
ExecStartPre=/bin/bash %h/services/caddy/configs/scripts/seed-caddyfile.sh
ExecStartPost=-/bin/bash %h/services/caddy/configs/scripts/export-root-ca.sh

[Install]
WantedBy=default.target
