# registry/caddy/quadlets/caddy.container
[Unit]
Description=Caddy reverse proxy
StartLimitIntervalSec=60
StartLimitBurst=3

[Container]
ContainerName=caddy
Image=docker.io/library/caddy:2-alpine
AutoUpdate=registry
Network=caddy.network
PublishPort=${SERVICE_PORT_HTTP}:${SERVICE_PORT_HTTP}
PublishPort=${SERVICE_PORT_HTTPS}:${SERVICE_PORT_HTTPS}
Volume=%h/.local/share/services/caddy/data:/data:Z,U
Volume=%h/.local/share/services/caddy/config-state:/config:Z,U
Volume=%h/.local/share/services/caddy/config:/etc/caddy:Z
Pull=missing
EnvironmentFile=%h/.local/share/services/caddy/.env
HealthCmd=caddy version || exit 1

[Service]
EnvironmentFile=%h/.local/share/services/caddy/.env
Restart=always
RestartSec=5
TimeoutStartSec=300
ExecStartPre=/bin/bash %h/.local/share/services/caddy/configs/scripts/seed-caddyfile.sh
ExecStartPost=-/bin/bash %h/.local/share/services/caddy/configs/scripts/export-root-ca.sh

[Install]
WantedBy=default.target
