# registry/open-webui/quadlets/open-webui.container
[Unit]
Description=Open WebUI - AI chat interface
StartLimitIntervalSec=60
StartLimitBurst=3
After=network-online.target

[Container]
ContainerName=open-webui
Image=ghcr.io/open-webui/open-webui:0.8.12
Network=open-webui.network
PublishPort=${PORT_HTTP}:8080
Volume=open-webui-data.volume:/app/backend/data:U
EnvironmentFile=%h/services/open-webui/.env
Pull=missing
HealthCmd=curl -sf http://localhost:8080/ || exit 1

[Service]
EnvironmentFile=%h/services/open-webui/.env
ExecStartPre=/bin/bash %h/services/open-webui/configs/scripts/merge-ca-bundle.sh
ExecStartPost=/bin/bash -c "podman cp %h/services/open-webui/ca-bundle.crt open-webui:/usr/local/lib/python3.11/site-packages/certifi/cacert.pem 2>/dev/null || true"
Restart=always
RestartSec=5
TimeoutStartSec=300

[Install]
WantedBy=default.target
