# registry/posthog/quadlets/posthog-plugins.container
[Unit]
Description=PostHog plugins (CDP)
StartLimitIntervalSec=300
StartLimitBurst=5
Requires=posthog-postgres.service posthog-redis.service posthog-clickhouse.service posthog-kafka.service
After=posthog-postgres.service posthog-redis.service posthog-clickhouse.service posthog-kafka.service

[Container]
Image=docker.io/posthog/posthog:latest
ContainerName=posthog-plugins
Network=posthog.network:alias=plugins
Exec=./bin/plugin-server --no-restart-loop
Pull=never
EnvironmentFile=%h/.local/share/ryra/posthog/.env

[Service]
EnvironmentFile=%h/.local/share/ryra/posthog/.env
Restart=always
RestartSec=10
TimeoutStartSec=300

[Install]
WantedBy=default.target
