# crates/ryra-core/registry/authelia/quadlets/authelia.container
[Unit]
Description=Authelia SSO server
# Disable systemd's start rate limit. ryra restarts authelia on every
# `--auth` install/remove (client registration mutates configuration.yml),
# so a normal multi-service session can fire 4–5 restarts in under a
# minute. RestartSec=5 still rate-limits crash loops to 12/min.
StartLimitIntervalSec=0

[Container]
Image=docker.io/authelia/authelia:4.39
AutoUpdate=registry
ContainerName=authelia
Network=authelia.network
PublishPort=${SERVICE_PORT_HTTP}:9091
Volume=%h/.local/share/services/authelia/config:/config:Z
Pull=missing
EnvironmentFile=%h/.local/share/services/authelia/.env
HealthCmd=wget -q --spider http://localhost:9091/api/health || exit 1

[Service]
EnvironmentFile=%h/.local/share/services/authelia/.env
Restart=always
RestartSec=5
TimeoutStartSec=300
ExecStartPre=/bin/bash %h/.local/share/services/authelia/configs/scripts/generate-config.sh
ExecStartPre=/bin/bash %h/.local/share/services/authelia/configs/scripts/generate-users.sh
ExecStartPre=/bin/bash %h/.local/share/services/authelia/configs/scripts/generate-oidc-rsa-key.sh

[Install]
WantedBy=default.target
