# crates/ryra-core/registry/openclaw/quadlets/openclaw.container
[Unit]
Description=Openclaw
StartLimitIntervalSec=60
StartLimitBurst=3

[Container]
ContainerName=openclaw
Image=ghcr.io/openclaw/openclaw:2026.3.13-1
PublishPort=${SERVICE_PORT_GATEWAY}:18789
Volume=%h/.local/share/services/openclaw/config:/home/node/.openclaw:Z,U
Exec=sh -c "echo '{\"gateway\":{\"mode\":\"local\",\"bind\":\"lan\",\"auth\":{\"mode\":\"token\",\"token\":\"'\"$OPENCLAW_GATEWAY_TOKEN\"'\"},\"controlUi\":{\"dangerouslyAllowHostHeaderOriginFallback\":true}}}' > /home/node/.openclaw/openclaw.json && exec node openclaw.mjs gateway --allow-unconfigured"
Pull=missing
EnvironmentFile=%h/.local/share/services/openclaw/.env
HealthCmd=wget -q --spider http://localhost:18789/ || exit 1

[Service]
EnvironmentFile=%h/.local/share/services/openclaw/.env
Restart=always
RestartSec=5
TimeoutStartSec=300

[Install]
WantedBy=default.target
