# registry/vikunja/quadlets/vikunja.container
[Unit]
Description=Vikunja task management
StartLimitIntervalSec=60
StartLimitBurst=3
After=network-online.target

[Container]
ContainerName=vikunja
Image=docker.io/vikunja/vikunja:2.3.0
Network=vikunja.network
PublishPort=${PORT_HTTP}:3456
Volume=%h/.local/share/services/vikunja/db:/db:Z,U
Volume=%h/.local/share/services/vikunja/files:/app/vikunja/files:Z,U
EnvironmentFile=%h/.local/share/services/vikunja/.env
Pull=missing
HealthCmd=curl -sf http://localhost:3456/api/v1/info || exit 1

[Service]
EnvironmentFile=%h/.local/share/services/vikunja/.env
Restart=always
RestartSec=5
TimeoutStartSec=300

[Install]
WantedBy=default.target
