# crates/ryra-core/registry/twenty/quadlets/twenty.container
[Unit]
Description=Twenty CRM server
StartLimitIntervalSec=60
StartLimitBurst=3
Requires=twenty-postgres.service twenty-redis.service
# The worker consumes the email (and other bullmq) queues — without it,
# invites / password-reset mail land on the queue and sit there forever.
# Wants= so it's pulled in automatically without being failure-blocking.
Wants=twenty-worker.service
After=twenty-postgres.service twenty-redis.service

[Container]
Image=docker.io/twentycrm/twenty:v1.23.1
ContainerName=twenty
Network=twenty.network
PublishPort=${SERVICE_PORT_HTTP}:3000
Volume=%h/.local/share/services/twenty/data:/app/docker-data:Z
Pull=missing
EnvironmentFile=%h/.local/share/services/twenty/.env
HealthCmd=curl -sf http://localhost:3000/healthz || exit 1
HealthStartPeriod=60s
HealthInterval=30s
HealthRetries=5
HealthTimeout=5s

[Service]
EnvironmentFile=%h/.local/share/services/twenty/.env
Restart=always
RestartSec=5
TimeoutStartSec=600
ExecStartPre=/bin/bash -c 'for i in $(seq 1 30); do podman exec twenty-postgres pg_isready -U ${POSTGRES_USER:-twenty} && exit 0; sleep 2; done; exit 1'

[Install]
WantedBy=default.target
