# registry/forgejo/quadlets/forgejo-postgres.container
[Unit]
Description=Forgejo PostgreSQL database
StartLimitIntervalSec=60
StartLimitBurst=3

[Container]
Image=postgres:17-alpine
ContainerName=forgejo-postgres
Network=forgejo.network:alias=postgres
Volume=forgejo-db-data.volume:/var/lib/postgresql/data:U
Pull=missing
EnvironmentFile=%h/.local/share/ryra/forgejo/.env
HealthCmd=pg_isready -U ${POSTGRES_USER:-postgres} || exit 1

[Service]
EnvironmentFile=%h/.local/share/ryra/forgejo/.env
Restart=always
RestartSec=5
TimeoutStartSec=300

[Install]
WantedBy=default.target
