# registry/twenty/quadlets/twenty-worker.container
[Unit]
Description=Twenty CRM background worker
StartLimitIntervalSec=60
StartLimitBurst=3
Requires=twenty-postgres.service twenty-redis.service twenty.service
After=twenty-postgres.service twenty-redis.service twenty.service

[Container]
Image=docker.io/twentycrm/twenty:v1.23.1
ContainerName=twenty-worker
Network=twenty.network
Exec=yarn worker:prod
Volume=%h/.local/share/services/twenty/data:/app/docker-data:Z
Pull=missing
EnvironmentFile=%h/.local/share/services/twenty/.env
Environment=DISABLE_DB_MIGRATIONS=true

[Service]
EnvironmentFile=%h/.local/share/services/twenty/.env
Restart=always
RestartSec=5
TimeoutStartSec=300

[Install]
WantedBy=default.target
