# registry/twenty/quadlets/twenty-postgres.container
[Unit]
Description=Twenty PostgreSQL database
StartLimitIntervalSec=60
StartLimitBurst=3

[Container]
Image=docker.io/library/postgres:16-alpine
AutoUpdate=registry
ContainerName=twenty-postgres
Network=twenty.network:alias=db
Volume=twenty-db-data.volume:/var/lib/postgresql/data:U
Pull=missing
EnvironmentFile=%h/services/twenty/.env
HealthCmd=pg_isready -U ${POSTGRES_USER:-twenty} || exit 1

[Service]
EnvironmentFile=%h/services/twenty/.env
Restart=always
RestartSec=5
TimeoutStartSec=300

[Install]
WantedBy=default.target
