# registry/twenty/quadlets/twenty-worker.container
[Unit]
Description=Twenty CRM background worker
StartLimitIntervalSec=60
StartLimitBurst=3
Requires=twenty-postgres.service twenty-redis.service twenty.service
After=twenty-postgres.service twenty-redis.service twenty.service

[Container]
Image=twentycrm/twenty:latest
ContainerName=twenty-worker
Network=twenty.network
Exec=yarn worker:prod
Volume=twenty-data.volume:/app/docker-data:Z
Pull=missing
EnvironmentFile=%h/.local/share/ryra/twenty/.env
Environment=DISABLE_DB_MIGRATIONS=true

[Service]
EnvironmentFile=%h/.local/share/ryra/twenty/.env
Restart=always
RestartSec=5
TimeoutStartSec=300

[Install]
WantedBy=default.target
