# registry/synapse/quadlets/synapse-db.container
[Unit]
Description=Synapse Postgres database
StartLimitIntervalSec=60
StartLimitBurst=3

[Container]
Image=docker.io/library/postgres:16-alpine
ContainerName=synapse-db
Network=synapse.network:alias=db
Volume=synapse-db.volume:/var/lib/postgresql/data:U
EnvironmentFile=%h/.local/share/ryra/synapse/.env
# Synapse recommends C-collation UTF8 to avoid index-ordering bugs. Quadlet
# parses Environment= by whitespace, so the value must be double-quoted.
Environment="POSTGRES_INITDB_ARGS=--encoding=UTF8 --lc-collate=C --lc-ctype=C"
HealthCmd=pg_isready -U synapse -d synapse
HealthStartPeriod=10s
HealthInterval=10s
HealthRetries=6
HealthTimeout=5s
Pull=missing

[Service]
EnvironmentFile=%h/.local/share/ryra/synapse/.env
Restart=always
RestartSec=5
TimeoutStartSec=300

[Install]
WantedBy=default.target
