# registry/immich/quadlets/immich.container
[Unit]
Description=Immich photo management server
StartLimitIntervalSec=60
StartLimitBurst=3
Requires=immich-postgres.service immich-machine-learning.service immich-valkey.service
After=immich-postgres.service immich-machine-learning.service immich-valkey.service

[Container]
Image=ghcr.io/immich-app/immich-server:v2.7.4
ContainerName=immich-server
Network=immich.network
PublishPort=${RYRA_PORT_HTTP}:2283
Volume=%h/.local/share/ryra/immich/upload:/data:Z
Pull=missing
EnvironmentFile=%h/.local/share/ryra/immich/.env
HealthCmd=curl -sf http://localhost:2283/api/server/ping || exit 1

[Service]
EnvironmentFile=%h/.local/share/ryra/immich/.env
Restart=always
RestartSec=5
TimeoutStartSec=600
ExecStartPost=/bin/bash %h/.local/share/ryra/immich/configs/scripts/create-admin.sh
ExecStartPost=/bin/bash %h/.local/share/ryra/immich/configs/scripts/configure-oidc.sh

[Install]
WantedBy=default.target
