# registry/docuseal/quadlets/docuseal.container
[Unit]
Description=DocuSeal document signing
StartLimitIntervalSec=60
StartLimitBurst=3
After=network-online.target

[Container]
Image=docker.io/docuseal/docuseal:2.4.4
ContainerName=docuseal
Network=docuseal.network
PublishPort=${PORT_HTTP}:3000
# DocuSeal stores SQLite db, uploaded PDFs, signed outputs, and a secret
# encryption key under /data. Bind mount so the data survives container
# replacement and is inspectable from the host.
Volume=%h/services/docuseal/data:/data:Z,U
EnvironmentFile=%h/services/docuseal/.env
HealthCmd=wget -qO- http://localhost:3000/up >/dev/null || exit 1
HealthStartPeriod=30s
HealthInterval=15s
HealthRetries=5
HealthTimeout=5s
Pull=missing

[Service]
EnvironmentFile=%h/services/docuseal/.env
Restart=always
RestartSec=5
TimeoutStartSec=300

[Install]
WantedBy=default.target
